# deploy/neondb.service
[Unit]
Description=NeonDB Game Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=neondb
Group=neondb
WorkingDirectory=/var/lib/neondb
ExecStart=/usr/local/bin/neondb start
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal
SyslogIdentifier=neondb

[Install]
WantedBy=multi-user.target
