# deploy/voltra.service
[Unit]
Description=Voltra Game Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=voltra
Group=voltra
WorkingDirectory=/var/lib/voltra
ExecStart=/usr/local/bin/voltra start
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voltra

[Install]
WantedBy=multi-user.target
